# scripts/systemd/minipost.service
[Unit]
Description=minipost (FastAPI) Service
After=network.target

[Service]
# 运行目录
WorkingDirectory=/opt/minipost

# 读取环境变量（可选，不存在也不报错）
EnvironmentFile=/opt/minipost/.deploy.env
# 兜底默认
Environment=HOST=0.0.0.0
Environment=PORT=8000

# 关键：用 /bin/sh -lc 包一层，才能展开 ${HOST}/${PORT}
ExecStart=/bin/sh -lc '/opt/minipost/.venv/bin/uvicorn app.main:app --host "${HOST:-0.0.0.0}" --port "${PORT:-8000}"'

Restart=always
RestartSec=2
User=root
Group=root

# 直接打到 journal（用 journalctl 查看）
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
