# deploy/huandan.service
[Unit]
Description=Huandan Server (FastAPI)
After=network.target

[Service]
Environment=HUANDAN_BASE=/opt/huandan-server
Environment=HUANDAN_DATA=/opt/huandan-data
Environment=PORT=8000
Environment=HOST=0.0.0.0
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/opt/huandan-server
ExecStart=/opt/huandan-server/.venv/bin/python /opt/huandan-server/run.py
Restart=always
User=root
# 如需降低权限，可创建专用用户并调整目录权限

[Install]
WantedBy=multi-user.target
